# systemd/celestex-admin.service
[Unit]
Description=CELESTE X — Admin DB (SQLAdmin) on port 8000
After=network.target celestex.service

[Service]
Type=simple
User=celeste
Group=celeste
WorkingDirectory=/opt/celestex
Environment=PATH=/opt/celestex/.venv/bin
Environment=CELESTEX_DB_PATH=/opt/celestex/data/celestex.db
EnvironmentFile=-/opt/celestex/.env
ExecStart=/opt/celestex/.venv/bin/uvicorn backend_admin.main:app --host 0.0.0.0 --port 8000 --proxy-headers
Restart=always
RestartSec=5
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
